# supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:nginx]
command=/usr/sbin/nginx

[program:dynamic]
command=/usr/local/bin/start-dynamic-nginx.sh
